# deployment/cobl-nexus.timer
[Unit]
Description=Run cobl computeNexusExport every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
